# building-directory-install/vendor/elo-mt-usb/elo.service
[Unit]
Description=Invoke Elo service at system startup.

[Service]
Type=forking
ExecStart=/etc/opt/elo-mt-usb/loadEloMultiTouchUSB.sh

[Install]
WantedBy=graphical.target
